# images/centos_httpd_php:7_latest_7.3/supervisord.d/init-app.camtech.com.co.conf
[program:python-runServer]
command= bash -c 'source djangoenv/bin/activate && \
         prj1/manage.py runserver 0.0.0.0:8000  &'

directory=/var/www/html
stopasgroup=true

