# images/centos_node_react:6_latest_6.17.1/supervisord.d/init-app.camtech.com.co.conf
[program:react-runServer]
command= bash -c 'npm start'

directory=/var/www/html
stopasgroup=true


