# images/centos_django:1.7.4/supervisord.d/init-app.swapi.com.co.conf
[program:python-runServer]
command= bash -c 'source djangoenv/bin/activate && \
         swapi/manage.py runserver 0.0.0.0:8000  &'

directory=/var/www/html
stopasgroup=true

